# deployment/dmh-api.service
[Unit]
Description=DMH Digital Marketing Hub API Service
After=network.target mysql.service

[Service]
Type=simple
User=dmh
Group=dmh
WorkingDirectory=/opt/dmh/backend
Environment="APP_HOST=0.0.0.0"
Environment="APP_PORT=8889"
ExecStart=/opt/dmh/backend/dmh -f /opt/dmh/backend/api/etc/dmh-api.yaml
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dmh-api

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
